# dots/.config/autostart/blueman.desktop
[Desktop Entry]
Exec=blueman-applet
Hidden=true
